# net-print/iprint/files/novell-iprint-listener-gui
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# iprint-listener-gui daemon providing network printing

command="/opt/novell/iprint/bin/iprint-listener-gui"

depend() {
	need logger cupsd clock
}
