# net-print/iprint/files/novell-iprint-listener
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

# iprint-listener daemon providing network printing

depend() {
	need logger cupsd clock
}


command="/opt/novell/iprint/bin/iprint-listener"
